# etc/habstranslate.timer
[Unit]
Description=Run habstranslate every minute
Requires=tankbot.service

[Timer]
OnCalendar=minutely

[Install]
WantedBy=timers.target
